# adblock_dns.path
[Unit]
Wants=adblock_dns.service

[Path]
PathChanged=/etc/adblock_blacklist
PathChanged=/etc/adblock_whitelist

[Install]
WantedBy=default.target
