# adblock_dns.service
[Unit]
Description=Update named zones for ad- and malware- blocking

[Service]
Type=oneshot
ExecStart=/var/lib/named/generate-zonefile.sh
User=root
